# scripts/systemd/dcfagent.service.j2
[Unit]
Description=DCF Agent

[Service]
ExecStart=/tmp/dcfagent -onboarding-service-url {{ onboarder_url }}

[Install]
WantedBy=multi-user.target
